# deploy/systemd/akior-kiosk.service
# AKIOR Kiosk Service
# Runs Chromium in kiosk mode on VT7
#
# Install:
#   sudo cp deploy/systemd/akior-kiosk.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now akior-kiosk.service
#
# Recovery:
#   Ctrl+Alt+F7 to return to kiosk (always)
#   Ctrl+Alt+F3 for console/TTY access

[Unit]
Description=AKIOR Console Kiosk
Documentation=https://github.com/yosiwizman/jarvis-v5-os/blob/main/docs/KIOSK_SETUP.md
Documentation=https://github.com/yosiwizman/jarvis-v5-os/blob/main/docs/ops/KIOSK_TROUBLESHOOTING.md
After=network-online.target systemd-logind.service
Wants=network-online.target
Conflicts=getty@tty7.service

# Prevent infinite restart storms if Xorg cannot start (e.g., GPU driver issue)
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=akior-kiosk
PAMName=login
StandardInput=tty
TTYPath=/dev/tty7
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes

# Run ExecStartPre/ExecStartPost as root (ExecStart still runs as User=)
PermissionsStartOnly=true
WorkingDirectory=/home/akior-kiosk
Environment=HOME=/home/akior-kiosk
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/akior-kiosk/.Xauthority
Environment=KIOSK_URL=https://akior.local/menu
SyslogIdentifier=akior-kiosk
StandardOutput=journal
StandardError=journal

# Pre-start cleanup (best-effort; runs as root to remove root-owned locks)
ExecStartPre=-/usr/bin/pkill -u akior-kiosk -x chromium
ExecStartPre=-/usr/bin/pkill -u akior-kiosk -x openbox
ExecStartPre=-/usr/bin/pkill -u akior-kiosk -x startx
ExecStartPre=-/usr/bin/pkill -u akior-kiosk -x xinit
ExecStartPre=-/usr/bin/pkill -u akior-kiosk -x Xorg
ExecStartPre=/bin/sh -lc 'set -x; echo "[akior-kiosk] preflight: cleaning stale X state"; if [ -f /tmp/.X0-lock ]; then pid="$(cat /tmp/.X0-lock 2>/dev/null || true)"; if echo "$pid" | grep -qE "^[0-9]+$"; then if ps -p "$pid" -o comm= 2>/dev/null | grep -q "^Xorg"; then echo "[akior-kiosk] preflight: killing stale Xorg pid=$pid"; kill -TERM "$pid" 2>/dev/null || true; sleep 1; kill -KILL "$pid" 2>/dev/null || true; fi; fi; fi; rm -f /tmp/.X0-lock /tmp/.X11-unix/X0 2>/dev/null || true; ls -la /tmp/.X0-lock /tmp/.X11-unix/X0 2>/dev/null || true; echo "[akior-kiosk] preflight: done"; exit 0'

# Start X server on VT7 with kiosk xinitrc
ExecStart=/usr/bin/startx /home/akior-kiosk/.xinitrc -- :0 vt7 -nolisten tcp -nocursor

# Ensure VT7 is active (non-fatal; chvt requires root)
ExecStartPost=-/usr/bin/chvt 7

Restart=on-failure
RestartSec=2
TimeoutStartSec=30
TimeoutStopSec=15
KillMode=control-group

[Install]
WantedBy=multi-user.target
